# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   

[supervisord]
logfile=/viirs/log/supervisord.log 
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
childlogdir=/viirs/log/
umask=022
pidfile=/tmp/supervisord.pid
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:sniff_update]
command=sniff_update
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=0
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=0
